# crontab
* * * * * root docker run seemly_app ./app.js audit >> /var/log/cron.log 2>&1
